# deploy/systemd/sandoghchi-bot.service
[Unit]

Description=Sandoghchi Telegram Bot


[Service]

WorkingDirectory=/opt/boursepilot

ExecStart=/opt/boursepilot/.venv/bin/python tools/run_telegram_bot.py

Restart=always


[Install]

WantedBy=multi-user.target

